# dashboard/systemd/bazaar-dashboard.service
[Unit]
Description=Bazaar Trading Desk Dashboard
After=default.target

[Service]
Type=simple
WorkingDirectory=/home/catabolic_solutions/.openclaw/workspace
ExecStart=/usr/bin/python3 /home/catabolic_solutions/.openclaw/workspace/dashboard/scripts/serve_dashboard.py
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
