# dashboard/systemd/bazaar-dashboard.service
[Unit]
Description=Bazaar Trading Desk Dashboard
After=default.target

[Service]
Type=simple
WorkingDirectory=/var/www/bazaar
ExecStart=/usr/bin/python3 /var/www/bazaar/dashboard/scripts/serve_dashboard.py
Restart=always
RestartSec=2
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
